# curator4-1283201167.conf
[program:curator4-1283201167]
command=python3 /home/curator4/curator4-1283201167/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1283201167/test.err.log
stdout_logfile=/home/curator4/curator4-1283201167/test.out.log